# lark-multi-agent.service
[Unit]
Description=Lark Multi-Agent - Multi-bot bridge for OpenClaw
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=haipw
WorkingDirectory=/home/haipw/lark-multi-agent
ExecStart=/usr/bin/node dist/index.js config.json
Restart=always
RestartSec=5
Environment=NODE_ENV=production

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lark-multi-agent

[Install]
WantedBy=multi-user.target
